# resources/services/restic-backup.timer
[Unit]
Description=Run restic-backup.service

[Install]
WantedBy=timers.target

[Timer]
OnCalendar=hourly
RandomizedDelaySec=5min
Persistent=true
