# resources/services/restic-backup.timer
[Unit]
Description=Run restic-backup.service

[Install]
WantedBy=timers.target

[Timer]
OnBootSec=5min
OnUnitInactiveSec=1h
RandomizedDelaySec=10min
