# systemd-units/restic-backup.timer
[Unit]
Description=Run restic-backup.service

[Install]
WantedBy=timers.target

[Timer]
OnBootSec=5m
OnUnitInactiveSec=1h
RandomizedDelaySec=5min
